# root/etc/uci-defaults/40_luci-theme-oasisic
#!/bin/sh
# luci-theme-oasisic — register theme in LuCI configuration
# Reference: OpenWrt luci-theme-bootstrap root/etc/uci-defaults/30_luci-theme-bootstrap

uci -q batch <<-EOF
  set luci.themes.oasisic=/luci-static/oasisic
  commit luci
EOF

exit 0